# scsi-customize/scsi-profile/files/init.d/closewireless
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/openssh/files/sshd.rc6,v 1.13 2004/07/15 00:07:55 agriffis Exp $

depend() {
	need coldplug
}

start() {
	ebegin "Starting close wireless"
		/sbin/rmmod ipw2200
	eend $?
}

stop() {
	ebegin "Stopping wireless close"
		echo ok
	eend $?
}
